# systemd/.config/systemd/user/battery-notify.service
[Unit]
Description=Run battery-notify

[Service]
ExecStart=/home/juan/.local/bin/services/battery-notify
